# templates/servarr-sso/generated/quadlets/authentik.container
[Container]
Name=authentik
Image=ghcr.io/goauthentik/server:REPLACE_WITH_TAG

Network=sso

Environment=AUTHENTIK_HOST=https://auth.owowatdis.cloud
Environment=AUTHENTIK_COOKIE_DOMAIN=.owowatdis.cloud

Volume=authentik-data:/data

Restart=always
