# templates/servarr-sso/generated/quadlets/jellyfin.container
[Container]
Name=jellyfin
Image=docker.io/jellyfin/jellyfin:REPLACE_WITH_TAG

Network=sso

Volume=jellyfin-config:/config
Volume=jellyfin-cache:/cache

Restart=always
