# templates/servarr-sso/generated/quadlets/lidarr.container
[Container]
Name=lidarr
Image=lscr.io/linuxserver/lidarr:REPLACE_WITH_TAG

Network=sso

Environment=PUID=1000
Environment=PGID=1000
Environment=TZ=UTC

Volume=lidarr-config:/config

Restart=always
