# templates/servarr-sso/generated/quadlets/navidrome.container
[Container]
Name=navidrome
Image=docker.io/deluan/navidrome:REPLACE_WITH_TAG

Network=sso

Volume=navidrome-data:/data

Restart=always
