# templates/servarr-sso/quadlets/jellyfin.container
[Container]
Name=jellyfin
Image=__JELLYFIN_IMAGE__

Network=sso

Volume=jellyfin-config:/config
Volume=jellyfin-cache:/cache

Restart=always
