# templates/servarr-sso/quadlets/kavita.container
[Container]
Name=kavita
Image=__KAVITA_IMAGE__

Network=sso

Volume=kavita-data:/kavita/config

Restart=always
