# templates/servarr-sso/quadlets/navidrome.container
[Container]
Name=navidrome
Image=__NAVIDROME_IMAGE__

Network=sso

Volume=navidrome-data:/data

Restart=always
