# templates/servarr-sso/quadlets/vaultwarden.container
[Container]
Name=vaultwarden
Image=__VAULTWARDEN_IMAGE__

Network=sso

Volume=vaultwarden-data:/data

Restart=always
